# docker/etc/supervisor.d/crontab.conf
[program:crontab]
command = coffee /data/htdocs/crontab.coffee
directory = /data/htdocs
stdout_logfile = /tmp/crontab.log
autostart = true
autorestart = true
