# docker/etc/supervisor.d/socket.io.conf
[program:socketio]
command = coffee /data/htdocs/socket.coffee
directory = /data/htdocs
stdout_logfile = /tmp/socket.io.log
autostart = true
autorestart = true
